# systemd/nft-auth-whitelist-server.service
# Example unit for the auth-server on the foreign (RFC) machine.
# This is a SAMPLE. Review it and the config before enabling. Do NOT expose
# Basic Auth over plain HTTP; put this behind HTTPS (Caddy/Nginx).
[Unit]
Description=nft-auth-whitelist auth-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/nft-auth-server --config /etc/nft-auth-whitelist/server.json
Restart=on-failure
RestartSec=3
# Hardening (review for your environment).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nft-auth-whitelist /var/log/nft-auth-whitelist
PrivateTmp=true

[Install]
WantedBy=multi-user.target
